# docker/supervisord.conf
# mathom — supervisord runs the two runtime processes side by side in one
# container: the FastAPI backend (uvicorn) and nginx. Both run as the
# non-root mathom user; all logs are forwarded to the container's stdout/stderr.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:backend]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/srv/mathom
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
